# nezha.ini
[program:nezha]
command=/usr/bin/nezha-agent -c /etc/nezha-agent/config.yml
directory=/
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nezha.err.log
stdout_logfile=/var/log/nezha.out.log
user=root
environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT="0"
