# node.ini
[program:node]
command=/usr/local/bin/node index.js
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/node.err.log
stdout_logfile=/var/log/node.out.log
user=root
environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT="0"

[program:node2]
command=/usr/local/bin/node ws.js
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/node.err.log
stdout_logfile=/var/log/node.out.log
user=root
environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT="0"
